# provisioning/marge/files/systemd/podman-network.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/podman-network.sh

[Install]
WantedBy=default.target
